# config/systemd/user/aria2c.service
[Unit]
Description=aria2c
Documentation=man:aria2c(1)
PartOf=graphical-session.target

[Service]
ExecStart=aria2c
