# config/systemd/user/dropbox.service
[Unit]
Description=dropbox
PartOf=graphical-session.target

[Service]
ExecStart=dropbox
