# config/systemd/user/greenclip.service
[Unit]
Description=greenclip
PartOf=graphical-session.target

[Service]
ExecStart=greenclip daemon
Restart=on-failure
