# config/systemd/user/lxpolkit.service
[Unit]
Description=lxpolkit
PartOf=graphical-session.target

[Service]
ExecStart=lxpolkit
Restart=on-failure
