# config/systemd/user/pa-applet.service
[Unit]
Description=pa-applet
Documentation=man:pa-applet(1)
PartOf=graphical-session.target

[Service]
ExecStart=pa-applet
