# config/systemd/user/plank.service
[Unit]
Description=plank
Documentation=man:plank(1)
PartOf=graphical-session.target

[Service]
ExecStart=plank
