# config/systemd/user/reload-systemd.path
[Unit]
Description=watch systemd config for changes

[Path]
PathModified=%h/.config/systemd/user/
