# config/systemd/user/ssh-agent.service
[Unit]
Description=SSH agent

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ssh-agent-per-user
RemainAfterExit=yes
