# config/systemd/user/syndaemon.service
[Unit]
Description=synaptic daemon
PartOf=graphical-session.target

[Service]
ExecStart=syndaemon -k
Restart=on-failure
