# config/systemd/user/udiskie.service
[Unit]
Description=udiskie
Documentation=man:udiskie(8)

[Service]
ExecStart=udiskie --tray
