# config/systemd/user/unclutter.service
[Unit]
Description=unclutter
Documentation=man:unclutter-xfixes(1)
PartOf=graphical-session.target

[Service]
ExecStart=unclutter
