# config/systemd/user/wallpaper.service
[Unit]
Description=wallpaper
Wants=autorandr.service
After=autorandr.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper
