# config/systemd/user/xscreensaver.service
[Unit]
Description=xscreensaver
PartOf=graphical-session.target

[Service]
ExecStart=xscreensaver
Restart=on-failure
